# deploy/supervisor/criserve-workers.conf
[program:criserve-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/criserve-portal/artisan queue:work redis --queue=default,notifications,mail --sleep=1 --tries=3 --timeout=180 --max-time=3600
directory=/var/www/criserve-portal
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/criserve-default.log
stopwaitsecs=3600

[program:criserve-documents]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/criserve-portal/artisan queue:work redis --queue=documents,exports --sleep=1 --tries=3 --timeout=1800 --max-time=3600
directory=/var/www/criserve-portal
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/criserve-documents.log
stopwaitsecs=3600

[program:criserve-imports]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/criserve-portal/artisan queue:work redis --queue=imports,deduplication,payouts --sleep=1 --tries=2 --timeout=3600 --max-time=3600
directory=/var/www/criserve-portal
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/criserve-imports.log
stopwaitsecs=3600

[program:criserve-audit]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/criserve-portal/artisan queue:work redis --queue=audit --sleep=1 --tries=3 --timeout=120 --max-time=3600
directory=/var/www/criserve-portal
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/criserve-audit.log
stopwaitsecs=3600
